# Linux/CentOS/config_bash_script_as_systemd_daemon_service_on_centos_7/data2report.service
[Unit]
Description=data2report

[Service]
Type=simple
ExecStart=/home/xxu/data2report/data2report

[Install]
WantedBy=multi-user.target

